# installer/macos/com.aetherav.realtime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.aetherav.realtime</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/aether</string>
        <string>watch</string>
        <string>/Users</string>
        <string>--quarantine</string>
        <string>/Library/Application Support/AetherAV/quarantine</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/AetherAV/realtime.log</string>
    <key>StandardOutPath</key>
    <string>/Library/Logs/AetherAV/realtime.log</string>
</dict>
</plist>
